# openrc/joycond
#!/sbin/openrc-run
depend() {
	need net
}
command="/usr/bin/joycond"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
